# files/supervisor-mongo.conf
[program:mongo]
autorestart=true
command=/usr/bin/mongod --config /etc/tokumx.conf
